# sup-uwsgi.ini
[program:sup-uwsgi]
command=uwsgi uwsgi.ini
autostart=true
autorestart=unexpected
startsecs=1
redirect_stderr=true
;stdout_logfile=/tmp/supervisor/uwsgi_stdout.log
